# src/azbridge/azbridge.service
[Unit]
Description=Azure Relay Bridge
  
[Service]
ExecStart=/usr/shared/azbridge
    
[Install]
WantedBy=multi-user.target